# etc/cron.d/router-reboot
MAILTO=andy@andyrobb.com
0 3 * * * pi /home/pi/git/lrt224-reboot/router-reboot.sh
